# tests/diverse-aur-test.service
[Unit]
Description=Diverse AUR Package Test v3
After=network.target

[Service]
Type=simple
User=deck
Environment=HOME=/home/deck
Environment=PATH=/home/deck/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/bin
ExecStart=/usr/bin/bash /tmp/diverse-aur-run-v3.sh
StandardOutput=append:/tmp/diverse-aur-results-v3.log
StandardError=append:/tmp/diverse-aur-results-v3.log
Restart=no

[Install]
WantedBy=multi-user.target
